# motd/gitkeeper-motd.service
[Unit]
Description=Generate GitKeeper MOTD
After=network.target

[Service]
Type=oneshot
ExecStart=gitkeeper motd
