# debian/git-keeper.timer
[Unit]
Description=Generate git-keeper MOTD hourly (timer)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
